# system-test/src/test/resources/mesos-es/slave-dind/supervisor.conf
[supervisord]
user=root
nodaemon=true

[program:docker]
user=root
command=/usr/local/bin/wrapdocker
redirect_stderr=true
stdout_logfile=/var/log/dind.log

[program:load]
user=root
command=bash -c 'while ! docker info; do sleep 1; done ; docker load < /tmp/build/images/executor.tar'
redirect_stderr=true
stdout_logfile=/var/log/docker-load.log

[program:mesos-slave]
user=root
command=bash -c 'while [ $(docker images | grep mesos/elasticsearch-executor) -eq "" ]; do sleep 1; done; mesos-slave'
redirect_stderr=true
stdout_logfile=/var/log/mesos-slave.log
