# system-test/src/test/resources/mesos-es/slave-dind/supervisor.conf
[supervisord]
user=root
nodaemon=true

[program:docker]
user=root
command=/usr/local/bin/wrapdocker
redirect_stderr=true
stdout_logfile=/var/log/dind.log

[program:import]
user=root
command=bash -c 'while ! docker info; do sleep 1; done ; docker import - mesos/elasticsearch-executor < /tmp/build/images/executor.tar'
redirect_stderr=true
stdout_logfile=/var/log/docker-import.log

[program:mesos-slave]
user=root
command=bash -c 'sleep 20; mesos-slave'
redirect_stderr=true
stdout_logfile=/var/log/mesos-slave.log
