# docker/dind/supervisor.conf
[supervisord]
user=root
nodaemon=true

[program:mesos-slave]
user=root
command=mesos-slave
redirect_stderr=true
stdout_logfile=/var/log/mesos-slave.log

[program:docker]
user=root
command=/usr/local/bin/wrapdocker
redirect_stderr=true
stdout_logfile=/var/log/dind.log
